# tomato-clock.service
# Tomato Clock - systemd service
# 1. 将本文件复制到 /etc/systemd/system/tomato-clock.service
# 2. 启动：sudo systemctl daemon-reload && sudo systemctl enable --now tomato-clock

[Unit]
Description=Tomato Clock PWA
After=network.target

[Service]
Type=simple
User=user1
WorkingDirectory=/home/user1/Claude_Code/Tomato_Clock
Environment=PORT=3000
Environment=HOST=0.0.0.0
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

# 资源限制（防止影响服务器其他服务）
MemoryMax=64M
CPUQuota=20%
TasksMax=32
LimitNOFILE=256

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
